# extras/pi_lib_systemd_system/speedSnitchWatcher.path
[Path]
PathModified=/home/pi/speedsnitch/speedsnitch

[Install]
WantedBy=default.target