# extras/pi_lib_systemd_system/speedSnitchWatcher.path
[Path]
PathModified=/opt/speedsnitch/speedsnitch

[Install]
WantedBy=default.target